# dev/mkosi/mkosi.extra/usr/lib/systemd/user/arlen-ai-agent.service
[Unit]
Description=Arlen AI agent (the background curation loop, user service)
Documentation=https://github.com/arlenos/arlen
# Subscribes to Event Bus triggers, runs the enabled behaviours (from
# /usr/share/arlen/agent/behaviours), and gates every proposed action (audited,
# fail-closed). Order it after the proxy + auditd it depends on (best-effort).
After=arlen-ai-proxy.service arlen-auditd.service
Wants=arlen-ai-proxy.service arlen-auditd.service

[Service]
Type=simple
ExecStart=/usr/lib/arlen/libexec/arlen-ai-agent
Restart=on-failure
RestartSec=2

# Point the system-daemon sockets (knowledge + event-bus, bound under /run/arlen)
# explicitly; the agent reads ARLEN_KNOWLEDGE_SOCKET with an ARLEN_DAEMON_SOCKET
# fallback, so set both. The audit ingest socket is per-user (its default).
Environment=ARLEN_KNOWLEDGE_SOCKET=/run/arlen/knowledge.sock
Environment=ARLEN_DAEMON_SOCKET=/run/arlen/knowledge.sock
Environment=ARLEN_CONSUMER_SOCKET=/run/arlen/event-bus-consumer.sock
# The config broker is a root SYSTEM service binding at /run/arlen (not this
# --user service's $XDG_RUNTIME_DIR). Point the agent at the system path so its
# executor_live / read-tier / autonomy come from the broker; absent the broker
# it falls back to ai.toml.
Environment=ARLEN_CONFIG_BROKER_SOCKET=/run/arlen/config-broker.sock

# Reads ai.toml + the baked behaviours (read-only); writes the durable undo log
# under ~/.local/state/arlen and connects to the per-user audit socket. Unix
# sockets + the session bus only (model traffic goes through the proxy).
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/state/arlen
ReadWritePaths=/run/user/%U/arlen
RestrictAddressFamilies=AF_UNIX
NoNewPrivileges=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
LockPersonality=yes

[Install]
WantedBy=default.target
