# dev/mkosi/mkosi.extra/usr/lib/systemd/user/arlen-dogfood.service
[Unit]
Description=Arlen in-VM KG-AI dogfood (one-shot loop driver, user service)
Documentation=https://github.com/arlenos/arlen
# Drives the KG-AI loop once per boot so the headless verify channel can grep the
# result off the serial journal: inject a file.opened, wait a promotion pass, then
# ask the AI daemon a question. Order it after the AI session daemons it drives
# (best-effort); llama-server load latency is absorbed by the tool's own retry.
After=arlen-ai-engine-daemon.service arlen-ai-proxy.service
Wants=arlen-ai-engine-daemon.service arlen-ai-proxy.service

[Service]
Type=oneshot
ExecStart=/usr/bin/arlen-dogfood
# The event-bus producer socket is a SYSTEM daemon under /run/arlen; the default
# matches but pin it so a path change here is explicit.
Environment=ARLEN_PRODUCER_SOCKET=/run/arlen/event-bus-producer.sock
# Markers go to the journal; the kernel cmdline forwards the journal to ttyS0, so
# DOGFOOD lines land in the serial log the verify channel reads.
StandardOutput=journal
StandardError=journal

# NO sandboxing on purpose: this is a dev test driver standing in for a normal
# user app (the harness) calling the AI daemon. The AI daemon resolves a caller's
# identity by reading /proc/<pid>/exe, which a same-uid reader can only do when
# the target is dumpable; systemd hardening (NoNewPrivileges + the Protect*/
# Restrict* set) clears dumpability, so a hardened dogfood is unresolvable and the
# query is refused. A real app caller is an ordinary unhardened process, so the
# dogfood must be too - confining it would test a configuration that never occurs.

[Install]
WantedBy=default.target
