# dev/mkosi/mkosi.extra/usr/lib/systemd/user/arlen-desktop-shell.service
[Unit]
Description=Arlen desktop shell (top bar, waypointer, overlays)
# Part of the graphical session: bound to it (stops when it stops) and ordered
# after it, so the compositor is up and has imported WAYLAND_DISPLAY into the
# --user environment before the shell launches. arlen-session starts the target
# once cosmic-comp is ready.
PartOf=graphical-session.target
After=graphical-session.target

[Service]
Type=simple
ExecStart=/usr/bin/arlen-desktop-shell
Restart=on-failure
RestartSec=2

[Install]
WantedBy=graphical-session.target
