# dev/mkosi/mkosi.extra/usr/lib/systemd/user/arlen-ai-daemon.service
[Unit]
Description=Arlen AI daemon (the conversational query layer, user service)
Documentation=https://github.com/arlenos/arlen
# Answers KG-grounded queries: reads the knowledge graph (capability-scoped) and
# forwards model completions through the ai-proxy. Order it after the proxy + the
# audit ledger it depends on (best-effort).
After=arlen-ai-proxy.service arlen-auditd.service
Wants=arlen-ai-proxy.service arlen-auditd.service

[Service]
Type=simple
ExecStart=/usr/bin/arlen-ai-daemon
Restart=on-failure
RestartSec=2

# The knowledge + event-bus daemons are SYSTEM services binding /run/arlen; both
# socket envs default to the per-user $XDG_RUNTIME_DIR, so point them at the
# system paths explicitly (the daemon honours these overrides first).
Environment=ARLEN_DAEMON_SOCKET=/run/arlen/knowledge.sock
Environment=ARLEN_CONSUMER_SOCKET=/run/arlen/event-bus-consumer.sock

# Reads ~/.config/arlen/ai.toml (read-only home is enough); talks Unix sockets +
# the session bus only (provider HTTP goes through the proxy, never direct).
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=/run/user/%U/arlen
RestrictAddressFamilies=AF_UNIX
NoNewPrivileges=yes
PrivateTmp=yes
ProtectKernelTunables=yes
ProtectKernelModules=yes
ProtectControlGroups=yes
RestrictNamespaces=yes
LockPersonality=yes

[Install]
WantedBy=default.target
