# dev/mkosi/mkosi.extra/usr/lib/systemd/user/arlen-powerd.service
[Unit]
Description=Arlen power daemon (user service)
Documentation=https://github.com/arlenos/arlen
# arlen-powerd aggregates UPower battery/AC/lid state into a coarse PowerState,
# publishes power.state on the event bus (push) and serves the latest snapshot
# over org.arlen.Power1 on the session bus (pull), so the shell battery
# indicator and the AI layer consume one shared value instead of each forking
# UPower. It reads UPower as a client on the system bus and talks only over
# AF_UNIX (the session bus, the event-bus producer socket and the per-user
# query socket it binds). Suspend/idle/profile management (PWR-R2..R7) lands
# here behind a capability gate later.
# Ordered after the audit daemon: powerd binds its query socket under the shared
# per-user runtime dir (/run/user/%U/arlen), and arlen-auditd is the Type=notify
# USER daemon that brings that dir up (it create_dir_all's it before binding its
# own sockets and only signals READY after). Ordering after arlen-event-bus does
# NOT work: the event bus is a SYSTEM unit, so a user-manager After= on it is void
# and powerd would race the dir, failing to bind its query socket under a
# not-yet-existent dir (query_socket::bind does not create the parent) until it
# recovers on a restart. The event-bus
# producer socket powerd publishes to is a system service already up before the
# user session, so no user-side ordering is needed for it.
After=arlen-auditd.service

[Service]
# Type=dbus: the service is "started" once it owns its bus name. It runs eagerly
# on boot (WantedBy=default.target) because it must aggregate UPower and publish
# power.state continuously, so clients reach the already-owned name; there is no
# D-Bus activation file (lazy start-on-first-use would defeat the continuous
# publish).
Type=dbus
BusName=org.arlen.Power1
# Production binary path (system-wide install). Dev runs invoke the debug binary
# directly; this file is the canonical production unit.
ExecStart=/usr/lib/arlen/libexec/arlen-powerd
# A failure to acquire the bus name or the producer socket exits non-zero;
# restart with back-off so a persistent fault surfaces rather than crash-looping.
Restart=on-failure
RestartSec=2

# powerd authenticates each privileged D-Bus caller (PWR-R5 profile-set, PWR-R2
# power actions) by resolving /proc/<caller-pid>/exe against the app-id allowlist
# (resolve_caller_app_id -> app_id_from_pid). A mount-namespace sandbox
# (ProtectSystem, ProtectHome, PrivateTmp, PrivateDevices, ProtectControlGroups,
# ProtectKernel*) makes /proc/<other-pid>/exe unreadable to a same-uid reader, so
# those options break that peer auth outright: the shell's power actions would be
# denied. This is the identical trap the audit daemon and consent broker document
# and omit for the same reason; they are therefore intentionally OMITTED here too.
# powerd's integrity rests on the bus-attested caller PID + the fail-closed app-id
# resolve, not on filesystem isolation, and it holds no persistent state (it only
# reads UPower and talks to logind / power-profiles-daemon over the system bus).
# Without ProtectSystem there is no read-only tree, so no ReadWritePaths is needed
# (arlen-auditd create_dir_all's the shared /run/user/%U/arlen the query socket
# binds under). The hardening that does NOT create a mount namespace is kept.
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
NoNewPrivileges=yes
UMask=0077

[Install]
WantedBy=default.target
