# dev/mkosi/mkosi.extra/usr/lib/systemd/user/arlen-auditd.service
[Unit]
Description=Arlen Audit Daemon (user service)
Documentation=https://github.com/arlenos/audit-daemon
# The daemon re-emits audit.* events on the Event Bus. This is a
# best-effort dependency, not a hard one: the daemon serves its
# sockets and writes the ledger whether or not the bus is up, so
# `After=` orders startup without `Requires=` blocking it.
After=arlen-event-bus.service

[Service]
# `Type=notify`: the daemon calls sd_notify(READY=1) once both the
# ingest and read sockets are bound, so dependent units (the AI
# layer) do not connect before the audit log is ready to record.
Type=notify
# Production binary path (system-wide install). Dev runs invoke the
# debug binary directly; this file is the canonical production unit.
ExecStart=/usr/lib/arlen/libexec/arlen-auditd
Restart=on-failure
RestartSec=2

# The audit daemon authenticates each ingest peer by resolving its
# /proc/<pid>/exe against the producer allowlist (SO_PEERCRED gives the
# PID; the exe is the identity). A mount-namespace sandbox (ProtectSystem,
# ProtectHome, PrivateTmp, PrivateDevices, ProtectControlGroups,
# ProtectKernel*) makes /proc/<other-pid>/exe unreadable to a same-uid
# reader, so those options would break peer auth outright (the daemon then
# resets every ingest and the AI layer's fail-closed audit refuses every
# query). They are therefore intentionally OMITTED. The hardening that does
# NOT create a mount namespace is kept (no-new-privs, address-family +
# namespace + personality + W^X restrictions); the ledger's integrity rests
# on its HMAC key + hash chain, not on filesystem isolation. Without
# ProtectSystem there is no ReadWritePaths to mismatch, so the daemon
# create_dir_all's its own per-user ledger dir on first run.
RestrictAddressFamilies=AF_UNIX
RestrictNamespaces=yes
LockPersonality=yes
MemoryDenyWriteExecute=yes
NoNewPrivileges=yes

[Install]
WantedBy=default.target
