# dev/mkosi/mkosi.extra/usr/lib/systemd/user/arlen-ai-engine-daemon.service
[Unit]
Description=Arlen AI engine daemon (pi-backed agent loop + gate/audit/compensation, user service)
Documentation=https://github.com/arlenos/arlen
# The pi-agent-adoption drop-in for the retired native ai-agent (pi-agent-
# adoption.md step 9). It owns org.arlen.AI1 (the ai-daemon replacement) +
# org.arlen.AIAgent1 (the pull-transparency + undo surface), spawns the confined
# pi sidecar (node >=22 + the pi engine, bwrap-jailed) only when [ai] enabled,
# and gates/audits/compensates every proposed action fail-closed. Order it after
# the proxy + auditd + config-broker it depends on (best-effort Wants).
After=arlen-ai-proxy.service arlen-auditd.service arlen-config-broker.service
Wants=arlen-ai-proxy.service arlen-auditd.service

[Service]
Type=simple
ExecStart=/usr/lib/arlen/libexec/arlen-ai-engine-daemon
Restart=on-failure
RestartSec=2

# System-daemon sockets (knowledge + event-bus, bound under /run/arlen). The
# engine reads ARLEN_DAEMON_SOCKET (knowledge) with an XDG-runtime fallback, so
# point it at the system path; the audit ingest socket is per-user (its default).
Environment=ARLEN_DAEMON_SOCKET=/run/arlen/knowledge.sock
Environment=ARLEN_KNOWLEDGE_SOCKET=/run/arlen/knowledge.sock
Environment=ARLEN_CONSUMER_SOCKET=/run/arlen/event-bus-consumer.sock
# The config broker is a root SYSTEM service binding at /run/arlen (not this
# --user service's $XDG_RUNTIME_DIR). Point the engine at the system path so its
# executor_live / read-tier / autonomy come from the broker; absent the broker it
# falls back to ai.toml.
Environment=ARLEN_CONFIG_BROKER_SOCKET=/run/arlen/config-broker.sock

# The pi sidecar layout baked into the image (decided deploy topology): the node
# >=22 runtime, the pi engine install, and the built Arlen pi-plugins extension
# entry. The sidecar reads all three; an unresolvable path leaves pi unspawned
# (the contract socket still serves fail-closed), so a wrong path is fail-safe,
# never a boot failure. The image build (dev/mkosi 05-ai.sh.chroot) must vendor
# node + pi + the built extension to these paths; that packaging step is the
# remaining half of the cutover.
Environment=ARLEN_PI_NODE_RUNTIME=/usr/lib/arlen/node-runtime
Environment=ARLEN_PI_INSTALL=/usr/lib/arlen/pi
Environment=ARLEN_PI_EXTENSION=/usr/lib/arlen/pi-plugins/dist/extension.js

# The daemon authenticates every engine-socket + rpc-proxy peer, and every
# compensate/undo D-Bus caller, by resolving the caller's /proc/<pid>/exe against
# the F3 app-id chain (SO_PEERCRED gives the pid; ConnectionAuth::extract_from ->
# app_id_from_pid reads the exe). A mount namespace on THIS daemon (created by ANY
# of ProtectSystem, ProtectHome, PrivateTmp, PrivateDevices, ProtectKernel*) makes
# /proc/<other-pid>/exe unreadable to a same-uid reader, so it would reject EVERY
# connection and refuse EVERY undo - the whole engine IPC surface (AI requests from
# the shell, undo from the harness) dead in deployment. Proven on this host: a
# readlink of /proc/<pid>/exe under any single one of those options fails EACCES
# for a same-uid target. They are therefore OMITTED, exactly as the audit daemon
# and the consent broker omit them for the identical peer-auth reason. This does
# NOT weaken the real security boundary: the untrusted pi/node sidecar is confined
# by bwrap, which creates its OWN user + mount + pid + ipc + uts + net + cgroup
# namespaces (--unshare-*) that a non-sandboxed parent only makes freer, not
# weaker; the daemon itself is trusted first-party code whose integrity rests on
# SO_PEERCRED attestation + the gate + the audit ledger, not on its own filesystem
# isolation. RestrictNamespaces stays wide enough for bwrap's --unshare-* set. The
# hardening that does NOT create a mount namespace is kept.
RestrictAddressFamilies=AF_UNIX
LockPersonality=yes
RestrictNamespaces=user mnt pid ipc uts net cgroup
NoNewPrivileges=yes

[Install]
WantedBy=default.target
