# dev/mkosi/mkosi.extra/usr/lib/systemd/system/greetd.service.d/10-arlen-vt.conf
# greetd owns vt1 (the virtio-gpu scanout) for cosmic-comp, but Debian statically
# enables getty@tty1 on the same VT (a preset, so NAutoVTs=0 does not stop it).
# Conflict it: starting greetd stops getty@tty1 and frees vt1, so cosmic-comp's
# udev/KMS backend can take DRM master + scanout there.
[Unit]
Conflicts=getty@tty1.service
After=getty@tty1.service
