# daemons/installd/dist/permission-helper.service
[Unit]
Description=Arlen Permission Helper
Documentation=https://github.com/arlenos/installd

[Service]
Type=dbus
BusName=org.arlen.PermissionHelper1
ExecStart=/usr/lib/arlen/permission-helper
User=root

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ReadWritePaths=/var/lib/arlen/permissions
ProtectHome=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
