# dev/mkosi/mkosi.extra/usr/lib/systemd/system/arlen-graph.service
[Unit]
Description=Arlen Graph Daemon
Documentation=https://github.com/arlenos/knowledge
After=arlen-event-bus.service

[Service]
Type=simple
ExecStart=/usr/bin/arlen-graph-daemon
Restart=on-failure
RestartSec=1

# Paths. The daemon reads ARLEN_GRAPH_PATH (not ARLEN_GRAPH_DB) and binds its
# read-query socket at knowledge.sock; the knowledge-mcp server connects to the
# same path, so the two units must agree on it.
Environment=ARLEN_DAEMON_SOCKET=/run/user/%U/arlen/knowledge.sock
Environment=ARLEN_GRAPH_PATH=%h/.local/share/arlen/graph
Environment=ARLEN_DB_PATH=%h/.local/share/arlen/events.db
# Read graph.toml from /usr/share/arlen-config: neither the /root nor the
# /home/arlen overlay reliably surfaced the file to this SYSTEM daemon at runtime
# ("graph.toml not found" even with ProtectHome=no), whereas /usr provably bakes
# (every binary/unit/behaviour lives there) and is readable under ProtectSystem=
# strict. XDG_CONFIG_HOME makes dirs::config_dir() resolve to it, so the daemon
# reads /usr/share/arlen-config/arlen/graph.toml (which sets the project watch dir).
Environment=XDG_CONFIG_HOME=/usr/share/arlen-config
# Debug the watch-config resolution: logs the exact graph.toml path checked
# (loaded vs not-found), so the serial confirms config_dir() resolves right.
Environment=RUST_LOG=info,arlen_graph_daemon::project::watch_config=debug

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
# ProtectHome=no, NOT read-only: this daemon's PURPOSE is reading the user's home
# (project detection, file.opened path promotion, the ~/.timeline FUSE mount). As a
# SYSTEM (root) service, ProtectHome=read-only left it unable to read /home/arlen at
# all - graph.toml under /home/arlen/.config was "not found" though the --user AI
# daemons read their ai.toml from the same dir. The capability-scoped query socket,
# not ProtectHome, is this daemon's access boundary, so reading the user's home is
# by-design. (This also unbreaks the documented ~/.timeline FUSE conflict below.)
ProtectHome=no
ReadWritePaths=%h/.local/share/arlen
ReadWritePaths=/run/user/%U/arlen
# DEPLOYMENT CAVEAT (Phase 10, needs on-host verification): the daemon mounts a
# FUSE timeline at ~/.timeline (DEFAULT_TIMELINE_MOUNT). That conflicts with the
# hardening above twice over: ProtectHome=read-only makes ~/.timeline read-only
# (the mount point cannot be created/mounted), AND ProtectSystem/ProtectHome put
# this unit in a private mount namespace, so a FUSE mount made here is not
# visible to the user's session. Making the timeline usable under this unit
# likely needs ~/.timeline in ReadWritePaths plus mount propagation
# (PrivateMounts=no / a shared MountFlags), or relaxing ProtectHome for the
# mount point. Set ARLEN_TIMELINE_MOUNT to a ReadWritePaths-covered location
# meanwhile. The graph/DB/socket paths above are correct; only the FUSE mount is
# affected.

[Install]
WantedBy=default.target
