# dev/mkosi/mkosi.extra/usr/lib/systemd/system/arlen-graph.service
[Unit]
Description=Arlen Graph Daemon
Documentation=https://github.com/arlenos/knowledge
After=arlen-event-bus.service

[Service]
Type=simple
ExecStart=/usr/bin/arlen-graph-daemon
Restart=on-failure
RestartSec=1

# Paths. The daemon reads ARLEN_GRAPH_PATH (not ARLEN_GRAPH_DB) and binds its
# read-query socket at knowledge.sock; the knowledge-mcp server connects to the
# same path, so the two units must agree on it.
Environment=ARLEN_DAEMON_SOCKET=/run/user/%U/arlen/knowledge.sock
Environment=ARLEN_GRAPH_PATH=%h/.local/share/arlen/graph
Environment=ARLEN_DB_PATH=%h/.local/share/arlen/events.db
# Read graph.toml from the USER config home, not /root: the mkosi.extra/root/ tree
# does not bake into the image (the /home/arlen overlay does - the AI daemons read
# its ai.toml), so the config lives at /home/arlen/.config/arlen/graph.toml, which
# this daemon (a system service, ProtectHome=read-only) can still read. XDG_CONFIG_
# HOME points dirs::config_dir() there; without it a root service falls back to
# ~/Projects etc. and the project watcher finds nothing.
Environment=XDG_CONFIG_HOME=/home/arlen/.config
# Debug the watch-config resolution: logs the exact graph.toml path checked
# (loaded vs not-found), so the serial confirms config_dir() resolves right.
Environment=RUST_LOG=info,arlen_graph_daemon::project::watch_config=debug

# Security hardening
NoNewPrivileges=true
ProtectSystem=strict
ProtectHome=read-only
ReadWritePaths=%h/.local/share/arlen
ReadWritePaths=/run/user/%U/arlen
# DEPLOYMENT CAVEAT (Phase 10, needs on-host verification): the daemon mounts a
# FUSE timeline at ~/.timeline (DEFAULT_TIMELINE_MOUNT). That conflicts with the
# hardening above twice over: ProtectHome=read-only makes ~/.timeline read-only
# (the mount point cannot be created/mounted), AND ProtectSystem/ProtectHome put
# this unit in a private mount namespace, so a FUSE mount made here is not
# visible to the user's session. Making the timeline usable under this unit
# likely needs ~/.timeline in ReadWritePaths plus mount propagation
# (PrivateMounts=no / a shared MountFlags), or relaxing ProtectHome for the
# mount point. Set ARLEN_TIMELINE_MOUNT to a ReadWritePaths-covered location
# meanwhile. The graph/DB/socket paths above are correct; only the FUSE mount is
# affected.

[Install]
WantedBy=default.target
