# dev/mkosi/mkosi.extra/usr/lib/systemd/system/arlen-graph.service.d/10-system-sut.conf
# SUT override: run the knowledge daemon as a SYSTEM service on /run/arlen (the
# single-user image model), replacing the dist unit's per-user topology
# (/run/user/%U, %h - which resolve to /run/user/0 + /root under a system
# manager and would fail). The true per-user-manager topology is a later product
# refinement. The KG/DB live under the StateDirectory; the timeline FUSE mount
# moves out of ~/.timeline (ProtectHome conflict) to a StateDirectory path.
[Service]
Environment=ARLEN_DAEMON_SOCKET=/run/arlen/knowledge.sock
Environment=ARLEN_CONSUMER_SOCKET=/run/arlen/event-bus-consumer.sock
Environment=ARLEN_GRAPH_PATH=/var/lib/arlen/graph
Environment=ARLEN_DB_PATH=/var/lib/arlen/events.db
Environment=ARLEN_TIMELINE_MOUNT=/var/lib/arlen/timeline
ProtectHome=yes
ReadWritePaths=
RuntimeDirectory=arlen
RuntimeDirectoryPreserve=yes
StateDirectory=arlen
